# selkies/rootfs/etc/systemd/user/dbus.service.d/export-bus.conf
[Service]
# Export the D-Bus socket path for desktop sessions
Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=%t/bus"
